# etc/init.d/dbus-openrc
#!/sbin/openrc-run

description="Starts the D-Bus system message bus"

command="/usr/bin/dbus-daemon"
command_args="--system"
command_user="messagebus"
supervisor="supervise-daemon"

depend() {
    need procfs sysfs
    provide dbus
}

start_pre() {
    checkpath -d -m 0755 -o "${command_user}" /run/dbus
}
