# etc/init.d/dbus-openrc
#!/sbin/openrc-run

description="Starts the D-Bus system message bus"

command="/usr/bin/dbus-daemon"
command_args="--system"
command_user="messagebus"
pidfile="/run/dbus/pid"
supervisor=supervise-daemon

depend() {
    need procfs sysfs
    provide dbus
}

start_pre() {
    # --- DEBUGGING: Print state to screen and pause ---
    echo "--- DEBUG: About to start D-Bus. Pausing for 10 seconds. ---"
    sleep 10

    echo "--- DEBUG: Listing /run directory ---"
    ls -laR /run
    echo "--- DEBUG: Pausing for 10 seconds. ---"
    sleep 10

    if [ -x /usr/bin/dbus-uuidgen ]; then
        /usr/bin/dbus-uuidgen --ensure
    fi

    checkpath -d -m 0755 -o "${command_user}" /run/dbus
}

start() {
    # --- DEBUGGING: Attempt to start the daemon and show any errors ---
    echo "--- DEBUG: Attempting to start dbus-daemon now... ---"
    # We remove --quiet and run in the foreground to see all output
    start-stop-daemon --start --pidfile "${pidfile}" \
        --user "${command_user}" \
        --exec "${command}" -- ${command_args}
    local status=$?
    
    echo "--- DEBUG: dbus-daemon exited with status code: ${status}. Pausing for 30 seconds. ---"
    sleep 30

    # This part will likely not be reached if the daemon fails, but is here just in case.
    if [ ! -S /run/dbus/system_bus_socket ]; then
        echo "--- DEBUG: Socket was NOT created. ---"
    else
        echo "--- DEBUG: Socket WAS created successfully. ---"
    fi
    echo "--- DEBUG: End of D-Bus script. Pausing for 10 seconds. ---"
    sleep 10

    # Intentionally fail so OpenRC doesn't continue and clear the screen
    return 1
}
